# build/templates/border0-webui.service
[Unit]
Description=Border0 WiFi Gateway Management Service

[Service]
StartLimitInterval=10
StartLimitBurst=10
ExecStart=/opt/border0/webui/webui

WorkingDirectory=/opt/border0/webui

Restart=always

RestartSec=10
EnvironmentFile=-/etc/sysconfig/border0-gw
Environment="HOME=/root"

[Install]
WantedBy=multi-user.target
