# build/templates/border0-webui.service
[Unit]
Description=Border0 WiFi Gateway Management Service
After=network.target

[Service]
Type=simple
StartLimitInterval=10
StartLimitBurst=10
ExecStart=/opt/border0/webui/venv/bin/python3 /opt/border0/webui/webui.py

WorkingDirectory=/opt/border0/webui

Restart=always

RestartSec=10
EnvironmentFile=-/etc/sysconfig/border0-webui

[Install]
WantedBy=multi-user.target
