# build/templates/border0-metrics.service
[Unit]
Description=Border0 Metrics Collector Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/border0/webui
ExecStart=/opt/border0/webui/venv/bin/python3 /opt/border0/webui/metrics_collector.py
Restart=always

[Install]
WantedBy=multi-user.target